# feeds/irorikon/luci-app-transip/root/etc/uci-defaults/80_transip
#!/bin/sh

touch /etc/config/transip

return 0